# meta-netservice/recipes-netservice/wifi-offload-manager/files/wifi-offload-manager.service
[Unit]
Description=WiFi Offload Manager — multi-path network infrastructure service
Documentation=file:///etc/netservice/path-policies.json
After=network-pre.target wpa_supplicant.service
Wants=network-pre.target
Requires=wpa_supplicant.service

[Service]
# Type=simple: daemon does not call sd_notify() yet.
# P6-T3: change to Type=notify and add WatchdogSec=30s once Phase 6 watchdog
#         integration is implemented in the daemon.
Type=simple
ExecStart=/usr/sbin/wifi-offload-manager --config /etc/netservice/path-policies.json
Restart=on-failure
RestartSec=5s

# Runtime directory for the control socket
RuntimeDirectory=netservice
RuntimeDirectoryMode=0755

# Hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/sys/fs/cgroup/net_cls /etc/netservice /run/netservice

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-offload-manager

[Install]
WantedBy=multi-user.target
